# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

radarr_command="mono --debug /app/radarr/bin/Radarr.exe -nobrowser -data=/config"
echo ":: Starting Radarr ::"
exec s6-setuidgid abc $radarr_command